# tools/build/docker/s6/s6-rc.d/redis/run
#!/bin/sh
exec /usr/bin/redis-server /root/lanraragi/tools/build/docker/redis.conf
