# tools/build/docker/s6/services.d/lanraragi/run
#!/bin/sh
cd /home/dezhao/lanraragi/
export HOME=/root
perl /home/dezhao/lanraragi/script/launcher.pl -f /home/dezhao/lanraragi/script/lanraragi
