# tools/build/docker/s6/services.d/lanraragi/run
#!/bin/sh
cd /root/lanraragi/
export HOME=/root
perl /root/lanraragi/script/launcher.pl -f /root/lanraragi/script/lanraragi
