# tools/build/docker/s6/services.d/redis/run
#!/bin/sh
/usr/bin/redis-server /root/lanraragi/tools/build/docker/redis.conf
