# tools/build/docker/s6/s6-rc.d/lanraragi/run
#!/bin/sh
cd /root/lanraragi/
export HOME=/root
exec s6-setuidgid root perl /root/lanraragi/script/launcher.pl -f /root/lanraragi/script/lanraragi
